# deploy/atlas-bot.service
[Unit]
Description=Atlas Trading Bot — daily pipeline scheduler
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/atlas-trading-bot
ExecStart=/opt/atlas-trading-bot/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=30

# Environment
EnvironmentFile=/opt/atlas-trading-bot/.env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/atlas-trading-bot/data /opt/atlas-trading-bot/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-bot

[Install]
WantedBy=multi-user.target
