# deploy/kukulkan-dashboard.service
[Unit]
Description=Kukulkan Dashboard — Streamlit UI
After=network.target kukulkan-bot.service
Wants=kukulkan-bot.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/kukulkan-trade
ExecStart=/opt/kukulkan-trade/.venv/bin/streamlit run src/dashboard/app.py \
    --server.port=8501 \
    --server.address=127.0.0.1 \
    --server.headless=true \
    --browser.gatherUsageStats=false
Restart=on-failure
RestartSec=10

# Environment
EnvironmentFile=/opt/kukulkan-trade/.env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/kukulkan-trade/data /opt/kukulkan-trade/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kukulkan-dashboard

[Install]
WantedBy=multi-user.target
