# deploy/kukulkan-bot.service
[Unit]
Description=Kukulkan Trading Bot — daily pipeline scheduler
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=kukulkan
Group=kukulkan
WorkingDirectory=/opt/kukulkan-trade
ExecStart=/opt/kukulkan-trade/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=30

# Environment
EnvironmentFile=/opt/kukulkan-trade/.env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/kukulkan-trade/data /opt/kukulkan-trade/logs /home/kukulkan /opt/kukulkan-improve

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kukulkan-bot

[Install]
WantedBy=multi-user.target
